# teste_post/etc/crontab/rodar_deletar.sh
#!/bin/sh
php /home/aindex/www/cubesat/teste_post/etc/crontab/deletar.php